# template/docker/mongodb/conf/supervisor.conf
[supervisord]
nodaemon=true
[inet_http_server]
port = *:9001
username = admin
password = admin#@!123

[rpcinterface:supervisor] 
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[include]
files = /etc/supervisor/conf.d/*.conf

[program:mongod]
command=/usr/bin/mongod --dbpath /data/db --smallfiles
autostart=true
user=root
startsecs=3
redirect_stderr=true
stdout_logfile=/var/log/mongod.log

