# template/docker/app/conf/supervisor.conf
[supervisord]
nodaemon=true
[inet_http_server]
port = *:9001
username = admin
password = admin#@!123

[rpcinterface:supervisor] 
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[include]
files = /etc/supervisor/conf.d/*.conf

[program:nodejs]
directory=/apps
command=nohup /usr/local/bin/npm run start &
stdout_logfile=/apps/appserver_stdout.log 
